# claude-code/rootfs/etc/s6-overlay/s6-rc.d/init-claude/run
#!/command/with-contenv bashio
# ==============================================================================
# init-claude: one-shot initialization service
# Runs once at add-on startup. Sets up persistent HOME, API key, and environment.
# ==============================================================================

bashio::log.info "Initializing Claude Code add-on..."

# ------------------------------------------------------------------------------
# Persistent HOME
# The Claude Code binary is installed to /root/.local at build time.
# We mirror it to /data/home/.local so updates and config persist across restarts.
# ------------------------------------------------------------------------------
if [ ! -d /data/home/.local ]; then
    bashio::log.info "First run: setting up persistent home directory..."
    mkdir -p /data/home
    cp -a /root/.local /data/home/.local
fi

# Inject HOME and PATH into S6 container environment (picked up by all longrun services)
printf "%s" "/data/home" > /var/run/s6/container_environment/HOME
printf "%s" "/data/home/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin" \
    > /var/run/s6/container_environment/PATH

export HOME=/data/home
export PATH="/data/home/.local/bin:${PATH}"

# ------------------------------------------------------------------------------
# Anthropic API Key
# Comes from add-on options (password-masked, stored encrypted by HA Supervisor).
# Never hardcoded — injected at runtime only.
# ------------------------------------------------------------------------------
api_key=$(bashio::config 'api_key')
if bashio::var.has_value "${api_key}"; then
    printf "%s" "${api_key}" > /var/run/s6/container_environment/ANTHROPIC_API_KEY
    bashio::log.info "Anthropic API key configured"
else
    bashio::log.warning "No API key set — run 'claude auth login' in the terminal to authenticate"
fi

# ------------------------------------------------------------------------------
# Model override (optional)
# ------------------------------------------------------------------------------
model=$(bashio::config 'model')
if bashio::var.has_value "${model}"; then
    printf "%s" "${model}" > /var/run/s6/container_environment/ANTHROPIC_MODEL
    bashio::log.info "Model override: ${model}"
fi

# ------------------------------------------------------------------------------
# Auto-update
# ------------------------------------------------------------------------------
if bashio::config.true 'auto_update'; then
    bashio::log.info "Checking for Claude Code updates..."
    /data/home/.local/bin/claude update 2>&1 \
        || bashio::log.warning "Claude Code update check failed (non-critical, continuing)"
else
    printf "1" > /var/run/s6/container_environment/DISABLE_AUTOUPDATER
    bashio::log.info "Auto-update disabled"
fi

# ------------------------------------------------------------------------------
# User-defined environment variables (format: KEY=VALUE)
# ------------------------------------------------------------------------------
for env_entry in $(bashio::config 'environment_vars'); do
    key="${env_entry%%=*}"
    val="${env_entry#*=}"
    if [ -n "${key}" ]; then
        printf "%s" "${val}" > "/var/run/s6/container_environment/${key}"
        bashio::log.info "Custom env set: ${key}"
    fi
done

# ------------------------------------------------------------------------------
# Static environment
# ------------------------------------------------------------------------------
printf "0"              > /var/run/s6/container_environment/USE_BUILTIN_RIPGREP
printf "xterm-256color" > /var/run/s6/container_environment/TERM
printf "C.UTF-8"        > /var/run/s6/container_environment/LANG

# ------------------------------------------------------------------------------
# Working directory and CLAUDE.md
# /data/workdir is the default CWD inside tmux. Claude starts here.
# /homeassistant is the HA config dir — preferred CWD if mounted.
# ------------------------------------------------------------------------------
mkdir -p /data/workdir
mkdir -p /data/home/.claude/skills

# Write default CLAUDE.md to /data/workdir (overwritten each start to stay current)
cp /usr/share/claude-ha/CLAUDE.md /data/workdir/CLAUDE.md

# Append user custom instructions if provided
custom=$(bashio::config 'custom_instructions')
if bashio::var.has_value "${custom}"; then
    printf "\n\n---\n\n# User Custom Instructions\n\n%s\n" "${custom}" \
        >> /data/workdir/CLAUDE.md
    bashio::log.info "Custom instructions appended to CLAUDE.md"
fi

# Also place CLAUDE.md in /homeassistant if it's mounted and writable
if [ -d /homeassistant ] && [ -w /homeassistant ]; then
    # Only create if it doesn't already exist (respect user's own CLAUDE.md)
    if [ ! -f /homeassistant/CLAUDE.md ]; then
        cp /data/workdir/CLAUDE.md /homeassistant/CLAUDE.md
        bashio::log.info "Created /homeassistant/CLAUDE.md with HA context"
    fi
fi

# ------------------------------------------------------------------------------
# Init commands (user-defined shell commands run at startup)
# Useful for: installing additional MCP servers, tools, etc.
# ------------------------------------------------------------------------------
for cmd in $(bashio::config 'init_commands'); do
    bashio::log.info "Running init command: ${cmd}"
    eval "${cmd}" 2>&1 \
        || bashio::log.warning "Init command failed (non-critical): ${cmd}"
done

bashio::log.info "Claude Code initialization complete"
