# claude-code/rootfs/etc/s6-overlay/s6-rc.d/claude-ttyd/run
#!/command/with-contenv bashio
# ==============================================================================
# claude-ttyd: longrun service
# Starts ttyd web terminal server. Wraps Claude Code in tmux for session
# persistence (browser disconnect does not kill the Claude session).
#
# Architecture:
#   ttyd (WebSocket/xterm.js) → tmux session "claude" → /usr/local/bin/start-claude
#
# ttyd binds to the hassio interface only (not exposed externally).
# HA Ingress proxies WebSocket connections to this port.
# ==============================================================================

declare ingress_port
ingress_port=$(bashio::addon.ingress_port)

bashio::log.info "Starting ttyd web terminal on port ${ingress_port}..."

# Prefer /homeassistant as working dir (HA config dir) — most useful default.
# Fall back to /data/workdir which always exists.
if [ -d /homeassistant ] && [ -w /homeassistant ]; then
    cd /homeassistant
else
    cd /data/workdir
fi

exec ttyd \
    --port "${ingress_port}" \
    --interface hassio \
    --writable \
    --ping-interval 30 \
    -t "fontSize=14" \
    -t "rendererType=webgl" \
    -t "disableLeaveAlert=true" \
    -t "scrollback=5000" \
    tmux -u new-session -A -s claude /usr/local/bin/start-claude
